# systemdfiles/flaskapp.service

[Unit]
Description=flask app service

[Service]
ExecStart=python3 /home/ec2-user/cj1905/project/app.py

[Install]
WantedBy=multi-user.target
